# rootfs/etc/init.d/addhosts
#!/sbin/openrc-run

depend() {
  need mount9p
}

start() {
  if [ -f /vol/state/hosts ]; then
    ebegin "adding hosts to /etc/hosts"
    cat /vol/state/hosts >> /etc/hosts
    eend $?
  fi
}
