# rootfs/etc/init.d/extrattys
#!/sbin/openrc-run

start() {
  for ttypath in /sys/firmware/qemu_fw_cfg/by_name/opt/ttytab/*; do
    ttyname=${ttypath##*/}
    ttycommand=$(cat $ttypath)
    echo "configuring $ttyname"
    echo "$ttyname::respawn:$ttycommand" >> /etc/inittab
  done

  kill -HUP 1
}
