# rootfs/etc/init.d/mount9p
#!/sbin/openrc-run

start() {
  for virtiofs in /sys/bus/virtio/drivers/9pnet_virtio/*; do
    [ -d "$virtiofs" ] || continue
    read tag < "$virtiofs/mount_tag"
    [ "$tag" = root ] && continue
    mkdir -p "/vol/$tag"

    echo "mount virtiofs $tag on /vol/$tag"
    mount -t 9p -o version=9p2000.L,trans=virtio "$tag" "/vol/$tag"
  done
}

stop() {
  for virtiofs in /sys/bus/virtio/drivers/9pnet_virtio/*; do
    [ -d "$virtiofs" ] || continue
    read tag < "$virtiofs/mount_tag"
    [ "$tag" = root ] && continue

    if findmnt --mountpoint "/vol/$tag" > /dev/null; then
      echo "unmount virtiofs $tag from /vol/$tag"
      umount "/vol/$tag"
    fi
  done
}
