# rootfs/etc/init.d/syslogd
#!/sbin/openrc-run

start() {
  ebegin "starting syslogd"
  syslogd -K
  eend $?
}
