# rootfs/etc/inittab
ttyS0::sysinit:/etc/init.d/sysinit
ttyS0::wait:/sbin/openrc default
tty1::respawn:/bin/bash
ttyS2::respawn:/bin/bash

